# pi/etc/update-motd.d/24-uptime
#!/bin/sh

UPTIME=$(uptime -p)

echo " * Uptime:         $UPTIME"
