# pi/etc/update-motd.d/26-user
#!/bin/sh

echo " * Users logged in: $(who | awk '{print $1}' | sort | wc -l)"
who
echo
